# docker/files/etc/profile.d/aliases.sh
alias ls='ls --color=auto'
alias ll='ls -lah'
alias l.='ls -d .* --color=auto'
