# docker/files/etc/profile.d/git-prompt-apply.sh
PS1='🐳 $PWD $(__git_ps1 "\[\e[0;32m\](%s)\[\e[0m\]")\$ '
